# config/systemd/user/mbsync.timer
[Unit]
Description=Mailbox sync timer

[Timer]
OnBootSec=1m
OnUnitActiveSec=3m
Unit=mbsync.service

[Install]
WantedBy=timers.target
