# config/systemd/user/mbsync.service
[Unit]
Description=Mailbox sync service

[Service]
Type=oneshot
ExecStart=/usr/bin/mbsync -Va
