# config/systemd/user/msmtp.service
[Unit]
Description=MSMTP Queue Service

[Service]
Type=oneshot
ExecStart=/usr/local/bin/msmtp-runqueue.sh
